# files/rootdir/var/spool/cron/crontabs/www-data
MAILTO=""
* * * * * php /var/www/html/groupoffice/cron.php /var/www/html/groupoffice/config.php
